# tron-lifecycle-manager/crontab
# START CRON JOB
6 * * * * python3 /app/cleanup.py > /proc/1/fd/1 2>/proc/1/fd/2
# END CRON JOB
